# nutrifriend.service
[Unit]
Description=NutriFriend AI Server
Documentation=https://github.com/your-repo/nutrifriend-ai
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/nutrifriend-ai-2
Environment=PYTHONIOENCODING=utf-8
ExecStart=/bin/bash /home/pi/nutrifriend-ai-2/start_pi.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security options
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
