# nutrifriend.service
[Unit]
Description=NutriFriend AI Health Assistant
After=network.target bluetooth.target

[Service]
Type=simple
User=os
WorkingDirectory=/home/os/Project-Os3/Project-Os
ExecStart=/bin/bash -c 'source backend/.venv/bin/activate && ./start_pi.sh'
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
